# .config/systemd/user/claude-overlay.service
[Unit]
Description=Claude streaming text overlay
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/claude-overlay
Restart=always
RestartSec=3
StartLimitIntervalSec=0

[Install]
WantedBy=graphical-session.target
