# .config/systemd/user/claude-hid.service
[Unit]
Description=Claude HID event listener (USB button → action)
After=claude-voice.service

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/.local/share/claude-hid/.venv/bin/python3 %h/.local/share/claude-hid/daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
