# .config/systemd/user/gcal-refresh.service
[Unit]
Description=Refresh Google Calendar cache

[Service]
Type=oneshot
ExecStart=%h/.local/bin/gcal-refresh
