# .config/systemd/user/google-contacts-sync.service
[Unit]
Description=Sync Google Contacts

[Service]
Type=oneshot
ExecStart=%h/.local/bin/google-contacts-sync
