# .config/systemd/user/mic-monitor.service
[Unit]
Description=Mic monitor loopback (wireless mic -> headphones)

[Service]
ExecStart=%h/.local/bin/mic-monitor-start.sh
Restart=on-failure
