# .config/systemd/user/system-freshness.service
[Unit]
Description=Refresh system freshness cache

[Service]
Type=oneshot
ExecStart=%h/.local/bin/system-freshness --refresh
