# .config/systemd/user/voice-typing.service
[Unit]
Description=Voice typing daemon (Whisper + VAD)
After=pipewire-pulse.service
Requires=pipewire-pulse.service

[Service]
Type=simple
Environment=YDOTOOL_SOCKET=%t/.ydotool_socket
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/bash -c 'test -x %h/.local/share/voice-typing-linux/.venv/bin/python || (python3.11 -m venv %h/.local/share/voice-typing-linux/.venv && %h/.local/share/voice-typing-linux/.venv/bin/pip install -r %h/.local/share/voice-typing-linux/requirements.txt requests)'
ExecStart=%h/.local/share/voice-typing-linux/.venv/bin/python %h/.local/share/voice-typing-linux/enhanced-voice-typing.py
Restart=on-failure
RestartSec=3
# Give model loading time before systemd considers it stuck
TimeoutStartSec=120

[Install]
WantedBy=default.target
