# root/etc/udev/rules.d/50-usb-anker.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="291a", ATTR{idProduct}=="8346", RUN+="/bin/sh -c 'echo -1 > /sys$devpath/power/autosuspend'"
