# .config/systemd/user/aside-daemon.service.d/api-key.conf
[Service]
ExecStartPre=%h/.local/bin/aside-load-api-key
